# server_config/pypln-manager.conf
[program:pypln-manager]
command=pypln-manager
directory=/srv/pypln/project/
user=pypln
stdout_logfile=/srv/pypln/logs/pypln-manager.out
stderr_logfile=/srv/pypln/logs/pypln-manager.err
